# deploy/systemd-user/chimera-client.service
[Unit]
Description=CHIMERA Client (user service)
After=chimera-gateway.service
Requires=chimera-gateway.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__CHIMERA_ROOT__
ExecStart=/usr/bin/env bash -lc './scripts/chimera-runner.sh cli up --state-file docs/runtime_state_latest.json'
ExecStop=/usr/bin/env bash -lc './scripts/chimera-runner.sh cli down --state-file docs/runtime_state_latest.json'
StandardOutput=append:__CHIMERA_ROOT__/docs/chimera_client.service.log
StandardError=append:__CHIMERA_ROOT__/docs/chimera_client.service.log

[Install]
WantedBy=default.target
